# crontab
# Execute every day at 6AM
# The first "cd" command is here to go to the directory containing main.py.
0 6 * * * cd /app && python main.py > /dev/stdout
